# etc/cron.daily/wedro_chroot-cron.daily
#!/bin/sh

CUSTOM_VAR='/var/opt'
CHROOT_DIR="$CUSTOM_VAR/chroot"

check_mounted() {
    if [ -z "$(mount | grep $CHROOT_DIR)" ]; then
	echo "CHROOT sems unmounted. exiting"
	exit 1
    fi
}

#######################################################################

start() {
    check_mounted
    chroot $CHROOT_DIR run-parts /etc/cron.daily
}

#######################################################################

start
